# wipeon.conf
[program:wipeon]
command = python app.py
directory = /home/pi/wipeon2
autostart = true
autorestart = true
